# luci-app-homekit/files/etc/uci-defaults/40_luci-homekit
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@homekit[-1]
	add ucitrack homekit
	set ucitrack.@homekit[-1].init=homekit
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
